# systemds/host/netdive-agent.service
[Unit]
Description=Netdive Agent
Wants=network-online.target
After=network-online.target

[Service]
User=root
Restart=on-failure

#Change this line if you download ther

ExecStart=/usr/share/ablestack/ablestack-netdive/netdive/ablestack-netdive-1.00 agent \
--conf /usr/share/ablestack/ablestack-netdive/netdive/netdive.yml

[Install]
WantedBy=multi-user.target
